# etc/cron.d/statusengine
0 0 * * *	www-data	/opt/statusengine/cakephp/app/Console/cake cronjob -q
